# etc/sv/localizer/run
#!/bin/sh

sleep 60
exec /usr/sbin/localize --min-gid 1000 --min-uid 1000
